# qn-code-assistant.service
# QN Code Assistant - Systemd Service
#
# To install:
#   sudo cp qn-code-assistant.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable qn-code-assistant
#   sudo systemctl start qn-code-assistant
#
# To view logs:
#   sudo journalctl -u qn-code-assistant -f

[Unit]
Description=QN Code Assistant
After=network.target

[Service]
Type=simple
User=claude
WorkingDirectory=/opt/claude-web
ExecStart=/usr/bin/python3 /opt/claude-web/app.py
Restart=on-failure
RestartSec=5
KillMode=process
Environment=PYTHONUNBUFFERED=1
Environment=CLAUDE_CODE_EXPERIMENTAL_AGENT_TEAMS=1

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/claude-web/sessions /opt/claude-web/backups /opt/claude-web/certs
ReadWritePaths=/opt/claude-web/config.json /opt/claude-web/usage.json
CapabilityBoundingSet=
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
